# etc/init.d/nohang-desktop-openrc
#!/sbin/openrc-run
# Copyright (c) 2007-2021 The Gentoo Foundation
# Distributed under the terms of the GNU General Public License v2

description="Starts the nohang low memory handler"
command="/usr/sbin/nohang"
command_args="--monitor --config /etc/nohang/nohang-desktop.conf"
pidfile="/run/nohang-desktop.pid"

# This tells OpenRC that the daemon will fork into the background
# and that OpenRC should create the PID file for it.
command_background="yes"

depend() {
	# This service should start after local filesystems are mounted.
	# It doesn't have hard dependencies, so 'after' is appropriate.
	after localmount
}
