# debian/postinst
#!/usr/bin/env bash

# -- Add service to runlevel.

rc-update add nohang-desktop-openrc default

# -- Update services list.

rc-update -u
